# itamae/cookbooks/configurator/files/etc/cron.d/djbu-configurator
*/3 * * * * dj /usr/bin/djbu-pull-config 2>&1 | logger -t djbu-pull-config
