# itamae/roles/gjcase/files/etc/systemd/system/gjcase.service
[Unit]
Description=gjcase unicorn server
 
[Service]
Type=simple
SyslogIdentifier=gjcase-unicorn
User=dj
PIDFile=/tmp/unicorn-gjcase.pid
WorkingDirectory=/home/dj/app/gjcase/current
 
ExecStart=/home/dj/app/gjcase/shared/env.sh bundle exec unicorn -c /home/dj/app/gjcase/current/config/unicorn.conf
ExecReload=/bin/kill -s USR2 $MAINPID
ExecStop=/bin/kill -s QUIT $MAINPID
 
[Install]
WantedBy=multi-user.target
